# scripts/com.slickwebhook.email.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- 서비스 식별자 -->
    <key>Label</key>
    <string>com.slickwebhook.email</string>

    <!-- 실행 파일 경로 -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/email-monitor</string>
    </array>

    <!-- 환경변수 (config.email.ini에서 로드하지만 백업용) -->
    <key>EnvironmentVariables</key>
    <dict>
        <!-- Email 설정 -->
        <key>EMAIL_IMAP_SERVER</key>
        <string>imap.gmail.com:993</string>
        <key>EMAIL_USERNAME</key>
        <string>your-email@gmail.com</string>
        <key>EMAIL_PASSWORD</key>
        <string>your-app-password</string>
        <key>EMAIL_FOLDER</key>
        <string>INBOX</string>
        <key>EMAIL_POLL_INTERVAL</key>
        <string>30s</string>

        <!-- ClickUp 설정 -->
        <key>CLICKUP_API_TOKEN</key>
        <string>pk_your_token</string>
        <key>CLICKUP_LIST_ID</key>
        <string>901413896178</string>

        <!-- Slack 설정 (알림용) -->
        <key>SLACK_BOT_TOKEN</key>
        <string>xoxb-your-bot-token</string>
        <key>SLACK_CHANNEL_ID</key>
        <string>C0A5ZTLNWA3</string>
    </dict>

    <!-- 로그인 시 자동 시작 -->
    <key>RunAtLoad</key>
    <true/>

    <!-- 종료 시 자동 재시작 -->
    <key>KeepAlive</key>
    <true/>

    <!-- 로그 파일 -->
    <key>StandardOutPath</key>
    <string>/Users/REPLACE_WITH_USERNAME/.slickwebhook/email.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/REPLACE_WITH_USERNAME/.slickwebhook/email.error.log</string>

    <!-- 작업 디렉토리 -->
    <key>WorkingDirectory</key>
    <string>/Users/REPLACE_WITH_USERNAME/.slickwebhook</string>
</dict>
</plist>
